# helper_scripts/hedwig_listener_dev.service
[Unit]
Description=Starts the DEV listener "Worker" which reaches out to workflow API.
After=network.target

[Service]
Type=simple
User=hedwig_dev
Group=hedwig_dev
Environment="HEDWIG_ENV=dev"
Environment="REQUESTS_CA_BUNDLE=/etc/pki/tls/certs/ca-bundle.crt"
Environment="PREFECT_API_URL=https://prefect2.hedwig-workflow-api.niaiddev.net/api"
Environment="IMOD_DIR=/opt/rml/imod"
WorkingDirectory=/gs1/home/hedwig_dev
ExecStart=/gs1/home/hedwig_dev/dev/bin/prefect worker start --pool workpool

[Install]
WantedBy=multi-user.target
